# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

[program:api-server]
command=python3 -m uvicorn 'api_server:app' --host=0.0.0.0 --port=8080
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=appuser
autostart=true
autorestart=true

[program:supercronic]
command=/usr/local/bin/supercronic /app/cronjobs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=appuser
autostart=true
autorestart=true
